# src/etc/rc.d/bgpd
#!/bin/sh
#
# $OpenBSD: src/lib/libssl/src/etc/rc.d/bgpd,v 1.1.1.1 2014/04/13 15:16:31 miod Exp $

daemon="/usr/sbin/bgpd"

. /etc/rc.d/rc.subr

pexp="bgpd: parent.*"

rc_cmd $1
