# src/etc/rc.d/cron
#!/bin/sh
#
# $OpenBSD: src/lib/libssl/src/etc/rc.d/cron,v 1.1.1.1 2014/04/13 15:16:31 miod Exp $

daemon="/usr/sbin/cron"

. /etc/rc.d/rc.subr

rc_cmd $1
