# src/etc/rc.d/inetd
#!/bin/sh
#
# $OpenBSD: src/lib/libssl/src/etc/rc.d/inetd,v 1.1.1.1 2014/04/13 15:16:32 miod Exp $

daemon="/usr/sbin/inetd"

. /etc/rc.d/rc.subr

rc_cmd $1
