# src/etc/rc.d/named
#!/bin/sh
#
# $OpenBSD: src/lib/libssl/src/etc/rc.d/named,v 1.1.1.1 2014/04/13 15:16:32 miod Exp $

daemon="/usr/sbin/named"

. /etc/rc.d/rc.subr

pexp="named: \[priv\]"

rc_cmd $1
