# src/etc/rc.d/ripd
#!/bin/sh
#
# $OpenBSD: src/lib/libssl/src/etc/rc.d/ripd,v 1.1.1.1 2014/04/13 15:16:32 miod Exp $

daemon="/usr/sbin/ripd"

. /etc/rc.d/rc.subr

pexp="ripd: parent.*"

rc_cmd $1
